# crontabs/root
# Every day at 03:00
0 3 * * * /home/docker/cron_autopostgresqlbackup.sh
